# files/etc/cron.d/named-cache
51 6 * * 7 root /opt/router/scripts/services/named-cache.sh &>/dev/null
